# scripts/screenshot-rename.timer
[Unit]
Description=Rename gnome-screenshots

[Timer]
Unit=screenshot-rename.service
OnBootSec=10s

[Install]
WantedBy=multi-user.target
