# scripts/screenshot-rename.service
[Unit]
Description=Rename gnome-screenshots

[Service]
Type=oneshot
User=lucero
ExecStart=/usr/bin/env bash -c ~/.myscripts/screenshot-watcher.sh
