# scripts/screenshot-rename.service
[Unit]
Description=Rename gnome-screenshots

[Service]
Type=oneshot
ExecStart=/bin/bash /home/lucero/.myscripts/screenshot-watcher.sh