# scripts/screenshot-rename.timer
[Unit]
Description=Rename gnome-screenshots

[Timer]
OnUnitActiveSec=3s
OnBootSec=60s

[Install]
WantedBy=timer.target